# util/docker/hpnp/service.full/hpnp.conf
[program:hpnp]
directory=/var/azuracast/scripts
command=/var/azuracast/scripts/hpnp.cjs
user=azuracast
priority=700
numprocs=1
