# util/docker/web/service.full/hpnp.conf
[program:hpnp]
directory=/var/azuracast
command=hpnp
user=azuracast
priority=700
numprocs=1
